# daemon/green-shark.service
[Unit]
Description=green-shark daemon

[Service]
ExecStart=green-shark

[Install]
WantedBy=default.target
